# opensmtpd/supervidord.ini
[supervisord]
nodaemon=true

[program:opensmtpd]
command=/usr/sbin/smtpd -d
autorestart=true
