# unbound/supervisord.ini
[supervisord]
nodaemon=true

[program:unbound]
command=/usr/sbin/unbound
autorestart=true
