# salt-minion/supervisord.ini
[supervisord]
nodaemon=true

[program:salt-minion]
command=/usr/bin/salt-minion
autorestart=true
